# systemd-hardening/systemd-rfkill.service.d/override.conf
[Service]
ProtectSystem=full
PrivateTmp=yes
NoNewPrivileges=yes
CapabilityBoundingSet=
