# systemd-hardening/emergency.service.d/override.conf
[Service]
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
