# systemd-hardening/plymouth-halt.service.d/override.conf
[Service]
NoNewPrivileges=yes
#ProtectSystem=full
#ProtectHome=yes
#PrivateTmp=yes
